# agent-bridge/ai.openclaw.poller.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.openclaw.poller</string>
    
    <key>Comment</key>
    <string>B.O.B. Message Poller - Checks Convex for messages every 30s</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>WorkingDirectory</key>
    <string>/Users/campoffice/.openclaw/workspace/agent-bridge</string>
    
    <key>ProgramArguments</key>
    <array>
      <string>/usr/local/bin/node</string>
      <string>/Users/campoffice/.openclaw/workspace/agent-bridge/poller.js</string>
    </array>
    
    <key>StandardOutPath</key>
    <string>/Users/campoffice/.openclaw/workspace/agent-bridge/poller.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/campoffice/.openclaw/workspace/agent-bridge/poller.err.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
      <key>HOME</key>
      <string>/Users/campoffice</string>
    </dict>
  </dict>
</plist>